# vendor/github.com/golang/gddo/deploy/services/redis/run
#!/bin/sh
exec /usr/bin/redis-server 2>&1
