# vendor/github.com/golang/gddo/deploy/services/gddo-server/run
#!/bin/sh
exec /go/bin/gddo-server 2>&1
